# deployments/launchd/com.paddock.prefetch-odds.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  締切前 live オッズ自動 prefetch（#237）を定期発火する launchd エージェント。
  __REPO_ROOT__ は install.sh が実リポジトリパスへ置換する（プレースホルダのまま load しない）。
  StartInterval=300（5 分間隔）。開催外/朝は対象 0 件で no-op になるため常時走らせてよい。
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.paddock.prefetch-odds</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__REPO_ROOT__/scripts/predict-check/prefetch_odds.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__REPO_ROOT__</string>

    <!-- launchd は最小 PATH なので psql / python3 / cargo 産物を探せるよう補う。
         Homebrew は Intel(/usr/local) と Apple Silicon(/opt/homebrew) の両方を入れておく。
         WORKDIR を固定し、launchd 起動でも本体ログが決まった場所（/tmp/paddock-prefetch/logs/）に
         出るようにする（README の確認パスと一致させる。$TMPDIR は launchd と対話シェルでズレうる）。 -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>WORKDIR</key>
        <string>/tmp/paddock-prefetch</string>
    </dict>

    <key>StartInterval</key>
    <integer>300</integer>

    <key>RunAtLoad</key>
    <true/>

    <!-- launchd の捕捉する stdout/stderr は常に存在する /tmp 直下へ出す（リポジトリ直下に
         未追跡ファイルを生成しない。本体の詳細ログは WORKDIR/logs/prefetch.log）。 -->
    <key>StandardOutPath</key>
    <string>/tmp/paddock-prefetch.launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/paddock-prefetch.launchd.err.log</string>
</dict>
</plist>
